# containers/compose/images/case5/platform/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
priority=100

[program:flask]
command=python3 /app/app.py
directory=/app
autorestart=true
priority=200